# 9/sys_check.timer
[Unit]
Description=Running sys_check service every 15 seconds

[Timer]
AccuracySec=1us
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
